# deploy/systemd/oir-mcp.service
[Unit]
Description=OIR MCP Server (AI Agent Interface)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/oir/repo
ExecStart=/opt/oir/venv/bin/python -m oir_mcp --transport streamable-http --port 8080
Restart=always
RestartSec=5
Environment=OIR_ROOT=/opt/oir
Environment=GITHUB_TOKEN_FILE=/opt/oir/.github-token

# Security hardening
NoNewPrivileges=yes
ReadWritePaths=/opt/oir
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oir-mcp

[Install]
WantedBy=multi-user.target
