# deploy/systemd/oir-webhook.service
[Unit]
Description=OIR GitHub Webhook Handler
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=oir
Group=oir
WorkingDirectory=/opt/oir/repo
ExecStart=/opt/oir/venv/bin/python /opt/oir/repo/deploy/webhook.py
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/oir
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oir-webhook

[Install]
WantedBy=multi-user.target
